# etc/profile.d/watchman-icad.csh
#!/bin/csh
setenv PATH /usr/local/watchman-icad/bin':'$PATH