# etc/profile.d/watchman-icad.sh
#!/bin/sh
PATH=/usr/local/watchman-icad/bin':'$PATH
export PATH